# ansible/templates/app.service.j2
[Unit]
Description=Sample App Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/app
ExecStart=/usr/bin/python3 /opt/app/server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
